# ansible/playbooks/iot/files/zigbee-soft-hang-watchdog.timer
[Unit]
Description=Run HAHA Zigbee soft-hang watchdog every 30s

[Timer]
# First run 60s after boot (let the stack settle), then every 30s.
OnBootSec=60s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=zigbee-soft-hang-watchdog.service

[Install]
WantedBy=timers.target
